# salt/freifunk/base/conntrack/etc/cron.d/conntrack
# (Salt managed)
# >/dev/null 2>/dev/null disables email alert sent by cron.d
0-59/5	*	* * * 	root /usr/local/bin/conntrack.sh >/dev/null 2>&1
