# salt/freifunk/base/ddmesh/etc/cron.d/freifunk-autoupdate
# (Salt managed)
# check Autoupdate every day
3 5 * * * root su -c "/usr/local/bin/freifunk-autoupdate >/dev/null 2>&1"
