# salt/freifunk/base/ddmesh/etc/cron.d/freifunk-masterless
### This file managed by Salt, do not edit by hand! ###
#
# Execute a local salt-call every 10 minutes
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/10 * * * * 	root /usr/bin/salt-call state.highstate --local >/dev/null 2>&1

### Execute after boot
@reboot			root /usr/bin/salt-call state.highstate --local >/dev/null 2>&1
