# salt/freifunk/base/bind/lib/systemd/system/bind9.service.tmpl
### This file managed by Salt, do not edit by hand! ###
{% from 'config.jinja' import nodeid %}
[Unit]
Description=BIND Domain Name Server
Documentation=man:named(8)
{# DNS Master or Slave Server #}
{% if nodeid == '3' or nodeid == '15' %}
After=network.target
{% else %}
After=network.target openvpn.service
{% endif %}

[Service]
ExecStart=/usr/sbin/named -f -u bind
ExecReload=/usr/sbin/rndc reload
ExecStop=/usr/sbin/rndc stop
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
