# salt/freifunk/base/cron/etc/cron.d/purge-old-kernels
# (Salt managed)
# purge-old-kernels and update grub every day

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

23 5 * * * root	/usr/bin/purge-old-kernels --keep 2 -qy && update-grub2 >/dev/null 2>&1
