# salt/freifunk/base/apt/etc/cron.d/apt-update
### This file managed by Salt, do not edit by hand! ###
3 */6 * * * root su -c "apt-get update > /dev/null 2>&1"
