# salt/freifunk/base/bind/slave/lib/systemd/system/bind9.service
### This file managed by Salt, do not edit by hand! ###
[Unit]
Description=BIND Domain Name Server
Documentation=man:named(8)
After=network.target

[Service]
ExecStart=/usr/sbin/named -f -u bind
ExecReload=/usr/sbin/rndc reload
ExecStop=/usr/sbin/rndc stop
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
